# local.d/11swapoff.start
#!/bin/bash
grep zram /proc/swaps | cut -d\  -f1 | xargs -r swapoff
