# local.d/69backlight.start
#!/bin/bash
[[ -e /sys/class/backlight/acpi_video0 ]] && echo 10 >/sys/class/backlight/acpi_video0/brightness
