# local.d/95custom-init.start
#!/bin/bash

[[ ! -x /home/init ]] || exec /home/init
